# project_v2/supervisor_files/gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/final_project/skillucket_v2
environment=SECRET_KEY="djangoinsecureoyjvqcqx6pogeuhx31b52qh7grv5mf1hp", DEBUG_SETTINGS=true
command=/home/ubuntu/final_project/skillucket_v2/env/bin/gunicorn --workers 3 --bind http://127.0.0.1:8000/ final_project.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[group:guni]
programs:gunicorn